# stow/linux/.config/systemd/user/kde-sycoca-rebuild.path
[Unit]
Description=Watch KDE service cache inputs

[Path]
PathChanged=%h/.config/mimeapps.list
PathChanged=%h/.local/share/applications
PathChanged=/usr/local/share/applications
PathChanged=/usr/share/applications
Unit=kde-sycoca-rebuild.service

[Install]
WantedBy=graphical-session.target
