# stow/linux/.config/systemd/user/kde-sycoca-rebuild.service
[Unit]
Description=Rebuild KDE service cache

[Service]
Type=oneshot
ExecStart=/usr/bin/kbuildsycoca6 --noincremental
Nice=5
