# stow/linux/.config/systemd/user/sunshine-audio.service
[Unit]
Description=Private PipeWire/Pulse sound server for Sunshine streaming (isolated from host audio)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.config/sway-sunshine/private-audio.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
