# stow/linux/.config/systemd/user/hyprwhspr.service.d/llm-postprocess.conf
[Service]
ExecStart=
ExecStart=/bin/bash -c 'export PYTHONPATH="$HOME/.config/hyprwhspr/patches:$PYTHONPATH"; exec /usr/lib/hyprwhspr/bin/hyprwhspr'
